# DEBIAN/prerm
#!/bin/sh

systemctl stop aqm-ip-display
systemctl disable aqm-ip-display

exit 0